# roles/raspy_config/templates/DrivingSwarmCore.service.j2
[Unit]
Description=driving swarm core
After=network-online.target

[Service]
Type=simple
User=turtle
Group=turtle
ExecStart=/home/turtle/core.sh

[Install]
WantedBy=multi-user.target
